# scripts/crontab
MAILTO=''
*/5 * * * *       cd /home/pi/wen && /usr/local/bin/rake code:update
